# dazaboost/ai.dazaboost.caffeinate.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>ai.dazaboost.caffeinate</string>

    <!-- Keeps the Mac awake ONLY during MLB game hours so the every-10-min poller doesn't miss
         pre-game windows. Starts 11:00am CT and holds awake for 13h (-t 46800) → ~midnight CT,
         then lets the machine sleep overnight.
         -i prevent idle sleep · -m prevent disk sleep · -s prevent system sleep (on AC). -->
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/caffeinate</string>
        <string>-i</string>
        <string>-m</string>
        <string>-s</string>
        <string>-t</string>
        <string>46800</string>
    </array>

    <!-- launch at 11:00am local (America/Chicago) each day; also start once when (re)loaded -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key><integer>11</integer>
        <key>Minute</key><integer>0</integer>
    </dict>
    <key>RunAtLoad</key>
    <true/>

    <key>StandardErrorPath</key>
    <string>/Users/jake/Dropbox/sbx/notebook/dazaboost/data/caffeinate.log</string>
</dict>
</plist>
